# systemd/redalert.service
[Unit]
Description=Red Alert Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
# Create QMI device node for bridge mode (qmicli requires a character device).
# The + prefix runs this command as root; the main process runs as User=pi.
ExecStartPre=+/bin/sh -c 'test -e /dev/cdc-wdm0 || mknod -m 0666 /dev/cdc-wdm0 c 180 0'
ExecStart=/home/pi/.local/bin/uv run python -m red_alert --config /home/pi/projects/red-alert/config.json
WorkingDirectory=/home/pi/projects/red-alert
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
